# test_em.timer
[Unit]
Description=test timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Unit=test_em.service

[Install]
WantedBy=timers.target